# crontab
0 0 * * * /scripts/backup.sh >> /backup/backup.log 2>&1
